# supervisor.d/postfix.ini
[program:postfix]
process_name	= postfix-monitor
directory	    = /etc/postfix
command     	= postfix.sh
startretries    = 0
startsecs	    = 0
autorestart	    = false
priority        = 100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[eventlistener:postfix_exit]
command=supervisor_killer postfix
process_name=postfix-monitor
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL,PROCESS_STATE_STOPPED,PROCESS_STATE_UNKNOWN
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
