# supervisor.d/rsyslog.ini
[program:rsyslog]
process_name	= syslog-runner
command     	= rsyslogd -n
startsecs	    = 0
autorestart	    = true
priority        = 0
user            = root
stdout_logfile  =/dev/stdout
stdout_logfile_maxbytes =0
stderr_logfile =/dev/stderr
stderr_logfile_maxbytes =0

[eventlistener:rsyslog_exit]
command=supervisor_killer syslog
process_name=syslog-runner
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL,PROCESS_STATE_STOPPED,PROCESS_STATE_UNKNOWN
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
